# docker/supervisord.conf

[program:nginx]
command=nginx -g "daemon off;"

[program:node]
command=node /var/www/api.localhost/dist/index.js